# botjam.crontab
@hourly /sbin/pfctl -t botjam -T expire 86400
